# bin/firstart.conf
start on login-session-start
script
    echo "Hey Joe!" >> /tmp/foobar.txt
    #cd /common/src/guadalinex/gecos-firstart/bin/
    #./firstart
end script
