# bin/firstart.conf
start on login-session-start
script
    echo "Hey Joe!" >> /tmp/foobar.txt
    cd /common/src/guadalinex/test-lock/bin
    ./firstart 2&>>/tmp/foobar.txt
end script
