# data/etc/xdg/autostart/firstart.desktop
[Desktop Entry]
Type=Application
Name=GECOS First-start assistant
Exec=/usr/bin/python __firstart_prefix__/bin/firstart
OnlyShowIn=GNOME;Unity;
NoDisplay=true
StartupNotify=false
X-GNOME-Autostart-Phase=Application
X-GNOME-Autostart-Notify=false
X-GNOME-AutoRestart=false
X-Ubuntu-Gettext-Domain=firstart
