# data/etc/init/firstart.conf
start on desktop-session-start
script
    __firstart_prefix__='/usr/local'
    cd "${__firstart_prefix__}/bin"
    ./firstart-dbusservice 2>/tmp/firstart-dbusservice.log
end script
