# data/etc/init/firstart.conf
start on desktop-session-start
script
    cd /usr/local/bin
    ./firstart-dbusservice 2>/tmp/firstart-dbusservice.log
end script
