# data/firstart.conf
start on desktop-session-start
script
    #sleep 10
    cd /usr/local/bin
    ./firstart-dbusservice 2>>/tmp/firstart-dbusservice.log
end script
